# deploy/streamlit.service
# systemd unit — DASHBOARD Streamlit (§5 dashboard público na VPS)
# Instalar na VPS:
#   sudo cp deploy/streamlit.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now streamlit
#   systemctl status streamlit
#
# O dashboard sobe na porta 8501. Exponha-a (firewall) ou coloque atrás de um
# proxy reverso (Nginx) para ter uma URL pública estável.
# AJUSTE os campos marcados com <...> ao seu servidor.
[Unit]
Description=Dashboard Streamlit (insights) — trabalho final
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<usuario_da_vps>
WorkingDirectory=<caminho_absoluto_do_repositorio>
EnvironmentFile=<caminho_absoluto_do_repositorio>/.env
ExecStart=<caminho_absoluto_do_repositorio>/.venv/bin/streamlit run app/streamlit_app.py \
    --server.port 8501 --server.address 0.0.0.0 --server.headless true
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
